# deploy/tenon.service
[Unit]
Description=Tenon woodworking app server
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=brian
ExecStart=/home/brian/.volta/bin/node /home/brian/current/server/dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=/etc/tenon/env

# DATA_DIR, PORT, MCP_BEARER_TOKEN are set in /etc/tenon/env
# Generate MCP_BEARER_TOKEN with: openssl rand -hex 32
# Set DATA_DIR=/home/brian/data (gitignored; backed up via Litestream + restic)

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenon

[Install]
WantedBy=multi-user.target
