# deploy/tenon.service
[Unit]
Description=Tenon woodworking app server
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=brian
ExecStart=/home/brian/.volta/bin/node /home/brian/current/server/index.js
Restart=always
RestartSec=5
EnvironmentFile=/etc/tenon/env

# /etc/tenon/env must define:
#   PORT=3000
#   DATA_DIR=/home/brian/data
#   MCP_BEARER_TOKEN=<openssl rand -hex 32>
#   NODE_ENV=production

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenon

[Install]
WantedBy=multi-user.target
